# Cygwin/etc/profile.d/esp_path.sh
#!/bin/sh
export IDF_PATH=/cygdrive/e/Aithinker/13_Aithinker_IDE/AiThinkerIDE_V0.5/cygwin/home/aithinker/project/esp-idf
export PATH=$PATH:/opt/xtensa-lx106-elf/bin:/opt/arm-none-eabi-gcc/bin:/opt/xtensa-esp32-elf/bin
